# raspberry-pi/sensor-node.service
[Unit]
Description=SmartHome BME690 MQTT sensor node
# Start after the network is up and, on the broker Pi, the local broker.
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
# Defaults target login user "raspberry" and ~/smarthome-sensor-node. On a Pi with a different user
# or path, run ./install.sh (it substitutes this Pi's user and directory) instead of hand-editing.
User=raspberry
WorkingDirectory=/home/raspberry/smarthome-sensor-node
EnvironmentFile=/home/raspberry/smarthome-sensor-node/sensor-node.env
ExecStart=/home/raspberry/smarthome-sensor-node/.venv/bin/python sensor_node.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
